# crontab
50 15 * * * root cd /app && /usr/local/bin/python price_change_detector.py >> /app/logs/fpl_price_changes.log 2>&1