# crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

45 08 * * * root /bin/bash -c 'cd /app && sleep 35 && /usr/local/bin/python price_change_detector.py'
